# deploy/signal-agent.service
[Unit]
Description=signal-agent daily healthcare digest (India + Global → Signal Agent India)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=signal
Group=signal
WorkingDirectory=/opt/signal-agent/repo
EnvironmentFile=/opt/signal-agent/shared/agent.env
ExecStart=/opt/signal-agent/repo/deploy/run-digest.sh india
# Build (~3 min) + hold-until-08:00 (up to ~90 min) + feedback/backup.
TimeoutStartSec=3600
Nice=5
